# lava-event-listener.service
[Unit]
Description=LAVA Device Health Event Listener
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lava-listener
Group=lava-listener
WorkingDirectory=/opt/lava-event-listener
ExecStart=/opt/lava-event-listener/.venv/bin/python -m lava_event_listener.main -c /etc/lava-event-listener/config.yaml
Restart=on-failure
RestartSec=10
StateDirectory=lava-event-listener

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/lava-event-listener
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
